# ata.path
[Path]
PathModified=/var/orbit/email_data/mail
Unit=ata.service
